# src/paccache.service.in
[Unit]
Description=Discard unused packages

[Service]
Type=oneshot
ExecStart=@bindir@/paccache -r
