# src/paccache.service.in
[Unit]
Description=Remove unused cached package files

[Service]
Type=oneshot
ExecStart=@bindir@/paccache -r
