# src/yaycache.service.in
[Unit]
Description=Remove unused cached package files

[Service]
Type=oneshot
ExecStart=@bindir@/yaycache -r
# Lowering priority
OOMScoreAdjust=1000
Nice=19
CPUSchedulingPolicy=idle
IOSchedulingClass=idle
IOSchedulingPriority=7
# Sandboxing and other hardening
NoNewPrivileges=yes
ProtectSystem=full
RestrictAddressFamilies=none
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallFilter=@system-service @file-system
SystemCallArchitectures=native
